# python-exporters/shelly/shelly1pm/shelly-mysql-exporter.service
[Unit]
Description=Shelly 1PM MySQL Exporter
After=network.target
 
[Service]
Restart=always
TimeoutStartSec=0
ExecStartPre=/bin/sleep 30
ExecStart=/usr/bin/python3 /root/python-exporters/shelly1pm2mysql.py
WorkingDirectory=/root/python-exporters
 
[Install]
WantedBy=multi-user.target
