# system-services/wakeonLAN.service
[Unit]
Description=Enable Wake On Lan

[Service]
Type=oneshot
ExecStart = /sbin/ethtool -s <interface HERE> wol g

[Install]
WantedBy=basic.target
